# docker/nginx/crontab
52 3,15 * * * root ssl-renew >> /dev/null
