# docker/nginx/crontab
52 3,15 * * * root /usr/local/bin/ssl-renew > /var/log/ssl-renew.log
